# deploy/systemd/gcal-organizer.service
[Unit]
Description=GCal Organizer — organize meeting docs, sync calendar, assign tasks
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=%h/.local/bin/gcal-organizer-wrapper.sh
EnvironmentFile=-%h/.gcal-organizer/.env
Environment=GCAL_DAYS_TO_LOOK_BACK=1

[Install]
WantedBy=default.target
