# deploy/systemd/gcal-organizer.timer
[Unit]
Description=GCal Organizer hourly timer

[Timer]
OnCalendar=hourly
Persistent=true
RandomizedDelaySec=120

[Install]
WantedBy=timers.target
