# deploy/launchd/com.jflowers.gcal-organizer.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.jflowers.gcal-organizer</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>WRAPPER_PATH_PLACEHOLDER</string>
    </array>

    <key>StartInterval</key>
    <integer>3600</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>LOG_PATH_PLACEHOLDER</string>

    <key>StandardErrorPath</key>
    <string>LOG_PATH_PLACEHOLDER</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
        <key>HOME</key>
        <string>HOME_PATH_PLACEHOLDER</string>
        <key>GCAL_ORGANIZER_BIN</key>
        <string>BINARY_PATH_PLACEHOLDER</string>
    </dict>
</dict>
</plist>
